# macos/headroom.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dotfiles.headroom</string>
  <key>ProgramArguments</key>
  <array>
    <!-- @HEADROOM@ is replaced with $(command -v headroom) by `make headroom`. -->
    <string>@HEADROOM@</string>
    <string>proxy</string>
  </array>
  <!-- Output shaping is a proxy-side feature, so it lives here (launchd won't
       inherit the shell's env). ANTHROPIC_BASE_URL is set in .zshrc instead —
       that one is read by Claude Code, not the proxy. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>HEADROOM_OUTPUT_SHAPER</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/headroom.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/headroom.log</string>
</dict>
</plist>
